# crontab
*/1 * * * * /var/www/temp.sh 2>/dev/null
*/3 * * * * /var/www/plot.sh > /var/www/graph.png
